# deploy/schedule-gunicorn.service
[Unit]
Description=Schedule App (Gunicorn)
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/apps/schedule/server
ExecStart=/home/pi/apps/schedule/server/.venv/bin/gunicorn -c gunicorn.conf.py wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
